# profiles/rmbp2018EndeavourOS/AutoExec.service
[Unit]
Description=AutoExec
[Service]
ExecStart=/home/hwt/.dotfiles/profiles/rmbp2018EndeavourOS/AutoExec.sh
[Install]
WantedBy=multi-user.target
